# deploy/systemd/gridrunner-events.service
[Unit]
Description=GRIDRUNNER event collection
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User={{OPERATOR_USER}}
WorkingDirectory={{OPERATOR_HOME}}
Environment=GRIDRUNNER_HOME={{GRIDRUNNER_HOME}}
Environment=GRIDRUNNER_OPERATOR_USER={{OPERATOR_USER}}
Environment=GRIDRUNNER_OPERATOR_HOME={{OPERATOR_HOME}}
Environment=GRIDRUNNER_EVENTS_LOG={{OPERATOR_HOME}}/{{OPERATOR_USER}}-events.log
EnvironmentFile=-{{GRIDRUNNER_HOME}}/state/gridrunner.env
TimeoutStartSec=90
SuccessExitStatus=124
ExecStart=/usr/bin/timeout 75s /usr/bin/bash {{OPERATOR_HOME}}/{{OPERATOR_USER}}-events.sh
